# configs/google-vrp/supervisor.conf
[supervisord]
nodaemon=true

[program:envoy-edge]
command=launch_envoy.sh -c /etc/envoy/envoy-edge.yaml %(ENV_ENVOY_EDGE_EXTRA_ARGS)s
  --log-format "(edge)[%%Y-%%m-%%d %%T.%%e][%%t][%%l][%%n] %%v" --base-id 0
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:envoy-origin]
command=launch_envoy.sh -c /etc/envoy/envoy-origin.yaml %(ENV_ENVOY_ORIGIN_EXTRA_ARGS)s
  --log-format "(origin)[%%Y-%%m-%%d %%T.%%e][%%t][%%l][%%n] %%v" --base-id 1
redirect_stderr=true
stdout_logfile_maxbytes=0
stdout_logfile=/dev/stdout
